# ops/systemd/governor-digest.service
[Unit]
Description=Nightly Memory Governor digest to markdown
After=memory-governor.service

[Service]
Type=oneshot
User=sacred
Group=sacred
WorkingDirectory=/opt/sacred-brain
Environment=MG_USER_ID=sam
Environment=DIGEST_OUTPUT_DIR=/opt/sam/memory/memory
Environment=MG_TZ_OFFSET_HOURS=11
ExecStart=/opt/sacred-brain/.venv/bin/python /opt/sacred-brain/scripts/governor_digest.py

# Hardening (lighter for oneshot)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/sam/memory
ReadOnlyPaths=/opt/sacred-brain
ProtectKernelTunables=true
ProtectKernelModules=true

[Install]
WantedBy=multi-user.target
