# ops/systemd/governor-digest.timer
[Unit]
Description=Nightly Memory Governor digest

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true

[Install]
WantedBy=timers.target
