# ops/systemd/hippocampus.service
[Unit]
Description=Sacred Brain Hippocampus Service
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/sacred-brain/hippocampus.env
ExecStart=/usr/local/bin/hippocampus
Restart=on-failure
User=sacred
Group=sacred

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/sacred-brain
ReadOnlyPaths=/opt/pipx /etc/sacred-brain
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
