# ops/systemd/hippocampus-memory-sync.timer
[Unit]
Description=Daily sync of Sam memory files into Hippocampus

[Timer]
OnCalendar=*-*-* 03:35:00
Persistent=true

[Install]
WantedBy=timers.target
