# ops/systemd/hippocampus-auto-tune.service
[Unit]
Description=Auto memory tuning for capture gating
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/sacred-brain
User=sacred
Group=sacred
Environment=HIPPOCAMPUS_SQLITE_PATH=/var/lib/sacred-brain/hippocampus/hippocampus_memories.sqlite
Environment=AUTO_TUNE_PATH=/var/lib/sacred-brain/auto_memory_tuning.json
ExecStart=/opt/pipx/venvs/sacred-brain-hippocampus/bin/python scripts/auto_memory_tuner.py

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/sacred-brain
ReadOnlyPaths=/opt/sacred-brain /opt/pipx
ProtectKernelTunables=true
ProtectKernelModules=true

[Install]
WantedBy=multi-user.target
