# ops/systemd/hippocampus.service
[Unit]
Description=Sacred Brain Hippocampus Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/sacred-brain
EnvironmentFile=/etc/sacred-brain/hippocampus.env
ExecStart=/opt/sacred-brain/.venv/bin/uvicorn brain.hippocampus.app:app --host 0.0.0.0 --port 54321
Restart=on-failure
User=sacred
Group=sacred

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/sacred-brain
ReadOnlyPaths=/opt/sacred-brain /etc/sacred-brain
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
