# ops/systemd/hippocampus-memory-sync.service
[Unit]
Description=Sync markdown memory files into Hippocampus
After=hippocampus.service

[Service]
Type=oneshot
User=sacred
Group=sacred
EnvironmentFile=/etc/sacred-brain/hippocampus.env
Environment=MEMORY_SYNC_ROOT=/opt/sam
Environment=HIPPOCAMPUS_USER_ID=sam
Environment=HIPPOCAMPUS_SQLITE_PATH=/var/lib/sacred-brain/hippocampus/hippocampus_memories.sqlite
ExecStart=/opt/sacred-brain/.venv/bin/python /opt/sacred-brain/scripts/memory_sync.py push

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/sacred-brain
ReadOnlyPaths=/opt/sacred-brain /opt/sam /etc/sacred-brain
ProtectKernelTunables=true
ProtectKernelModules=true

[Install]
WantedBy=multi-user.target
