# ops/systemd/hippocampus-notes-import.timer
[Unit]
Description=Nightly import of Org/Denote notes into Hippocampus memories

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
