# ops/systemd/matrix-autojoin.service
[Unit]
Description=Sacred Brain Matrix Auto-Join Helper
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/sacred-brain
EnvironmentFile=/opt/sacred-brain/.env.matrix
ExecStart=/opt/sacred-brain/.venv/bin/python bots/matrix/autojoin_bot.py
User=sacred
Group=sacred
Restart=on-failure

[Install]
WantedBy=multi-user.target
