# ops/systemd/memory-governor-consolidate.service
[Unit]
Description=Memory Governor consolidation run
After=memory-governor.service

[Service]
Type=oneshot
User=sacred
Group=sacred
EnvironmentFile=/etc/sacred-brain/memory-governor.env
ExecStart=/bin/bash -lc '\
  SCOPES="${MG_CONSOLIDATE_SCOPES}"; \
  if [ -z "$SCOPES" ]; then \
    echo "No MG_CONSOLIDATE_SCOPES set; skipping"; \
    exit 0; \
  fi; \
  IFS="," read -ra ITEMS <<< "$SCOPES"; \
  for s in "${ITEMS[@]}"; do \
    kind="${s%%:*}"; id="${s#*:}"; \
    [ -z "$kind" ] && continue; \
    curl -s -X POST -H "Content-Type: application/json" \
      -d "{\"scope\":{\"kind\":\"${kind}\",\"id\":\"${id}\"},\"mode\":\"all\"}" \
      http://127.0.0.1:54323/consolidate >/dev/null || true; \
  done'

# Hardening (lighter for oneshot)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true

[Install]
WantedBy=multi-user.target
